# systemd/star-tracker-ble.service
[Unit]
Description=Star Tracker BLE peripheral (goto_ble.py)
# bluetooth.target brings up bluetoothd; NetworkManager is needed by the
# `start_ap` BLE command (nmcli) and for address enumeration.
After=bluetooth.target NetworkManager.service
Requires=bluetooth.target
Wants=NetworkManager.service

[Service]
Type=simple
# Adjust these two paths if you cloned somewhere else. The venv is
# expected to have been created with `--system-site-packages` so that
# the apt-installed picamera2 / libcamera are visible — see README.
WorkingDirectory=/home/acidhax/dev/star_tracker_ws
ExecStart=/home/acidhax/dev/star_tracker_ws/.venv/bin/python3 /home/acidhax/dev/star_tracker_ws/goto_ble.py --log INFO

Restart=on-failure
RestartSec=3

# BLE peripheral has to advertise; nmcli needs root for hotspot control.
# Tighten with CAP_NET_ADMIN+CAP_NET_RAW if you want to drop privileges
# (you'll also need to grant the user access to /dev/ttyACM* and i2c).
User=root

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
